# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/tvheadend/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts TVHeadend
# ==============================================================================
declare args

bashio::log.info "Starting TVHeadend..."

args=$(bashio::config 'args')
bashio::log.info "Additional args ${args}"

exec /usr/bin/tvheadend -u root -g root -c /config/tvheadend "$args" 
