# tvheadend/rootfs/etc/cron.d/webgrabplus
0 0 * * * /config/tvheadend/wg++/run.sh
